# skrepka_bot.service
[Service]
ExecStart=/root/skrepka/Skrepka_tg_antiTinder_bot/start.sh
Restart=always

[Install]
WantedBy=multi-user.target